# system/etc/udev/rules.d/00-pok3r.rules
ACTION=="add", ATTRS{idVendor}=="04d9", ATTRS{idProduct}=="0141", RUN+="/usr/bin/usb-keyboard-in_udev"
